# backend/etc/nightly-build.timer
[Unit]
Description=Build GISST CI Container Nightly

[Timer]
OnCalendar=*-*-* 07,19:00:00
Persistent=true

[Install]
WantedBy=timers.target
